# osquery.conf
{
  "schedule": {
    "crontabs": {
      "query": "SELECT * FROM crontab;",
      "interval": 5
    }
  }
}
